# config/grip.crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
15 0 * * * grip-redis-backup.sh
